# config/systemd/user/git-annex-webapp.service
[Unit]
Description=git-annex webapp daemon

[Service]
ExecStart=/usr/bin/git-annex webapp
Restart=on-failure

[Install]
WantedBy=default.target
