# config/systemd/user/stoppidgin.service
[Unit]
Description=Stop pidgin

[Service]
Type=oneshot
ExecStart=/bin/echo ''
